# deploy/zenos-memory.service
[Unit]
Description=Zenos Memory personal sidecar
After=network-online.target

[Service]
Type=notify
NotifyAccess=all
User=zenos-memory
Group=zenos-memory
WorkingDirectory=/opt/zenos-memory/current
LoadCredentialEncrypted=zenos-memory.env:/etc/credstore.encrypted/zenos-memory.env.cred
Environment=HOME=/var/lib/zenos-memory
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3091
Environment=ZENOS_MEMORY_STORAGE_MODE=drive-events
Environment=ZENOS_MEMORY_AUTH_NONCE_STORE=process
Environment=ZENOS_MEMORY_DB_PATH=/var/lib/zenos-memory/zenos-memory.sqlite
ExecStart=/usr/bin/node /opt/zenos-memory/current/scripts/load-env-and-start.mjs
Restart=on-failure
RestartSec=3
TimeoutStartSec=60
TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=mixed
UMask=0077
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
DevicePolicy=closed
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
ReadOnlyPaths=/opt/zenos-memory/current /etc/zenos-memory
ReadWritePaths=/var/lib/zenos-memory
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
StateDirectory=zenos-memory
StateDirectoryMode=0700
CacheDirectory=zenos-memory
CacheDirectoryMode=0700
LimitNOFILE=32768
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
